# daemon/fs/etc/supervisor.d/aria2.ini
[program:aria2]
command=/usr/bin/aria2c
directory=/root/
autorestart=true
user=root
